# setup_robotino/udev_rules/99-joystick.rules
KERNEL=="js0", SUBSYSTEM=="input" RUN+="/etc/init.d/joystick.sh start"
